# modules/screen/.cron
* * * * * /home/andreid/Work/andreidmt.dotfiles/new/modules/screen/scripts/cron_camerashot
