# modules/pipewire/services/pipewire-pulse/run
#!/bin/sh

exec pipewire-pulse 2>&1 \
  | sed 's/^/[pipewire-pulse] /'
