# modules/@1_core/services/usb-automount/run
#!/bin/sh

exec 2>&1                            # stderr → stdout (for this shell and children)
exec "$(dirname "$0")/usb-automount" # replace this shell with the `usb-automount` script
