# modules/@1_core/services/usb-automount/run
#!/bin/sh 

exec 2>&1
exec "$(dirname "$0")/usb-automount" -u "$USER" \
  | sed 's/^/[usb-automount] /' 
