# modules/pipewire/services/pipewire/run
#!/bin/sh

# Set the pipewire log level.
# See: https://docs.pipewire.org/page_daemon.html#sec_logging
export PIPEWIRE_DEBUG="2"

exec pipewire 2>&1 \
  | sed 's/^/[pipewire] /'
