# modules/pipewire/services/wireplumber/run
#!/bin/sh

export WIREPLUMBER_DEBUG=2

exec wireplumber 2>&1 \
  | sed 's/^/[wireplumber] /'

