# modules/pipewire/services/mpd/run
#!/bin/sh

exec mpd --no-daemon --stderr 2>&1 \
  | sed 's/^/[mpd] /'
